# target/linux/ar71xx/base-files/etc/uci-defaults/05_branding
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

[ -e /etc/config/branding ] && exit 0

touch /etc/config/branding

. /lib/functions/uci-defaults.sh
. /lib/ar71xx.sh

board=$(ar71xx_board_name)

case "$board" in
carambola2)
	ucidef_set_system_hostname "Carambola2"
	ucidef_set_wifi_ssid "Carambola2"
	;;

centipede)
	ucidef_set_system_hostname "Centipede"
	ucidef_set_wifi_ssid "Centipede"

	# Centipede is not internet gateway device
	# WAN zone is used for eth0 (DHCP client)
	# Change firwall rule to allow access from eth0
	uci set firewall.@zone[1].input='ACCEPT'
	;;

lima)
	ucidef_set_system_hostname "Lima"
	ucidef_set_wifi_ssid "Lima"
	;;

*)
	ucidef_set_system_hostname "OpenWRT"
	;;
esac

uci commit system
uci commit wireless
uci commit firewall

exit 0
