# target/linux/ar71xx/base-files/etc/uci-defaults/05_branding
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

[ -e /etc/config/branding ] && exit 0

touch /etc/config/branding

. /lib/functions/uci-defaults.sh
. /lib/ar71xx.sh

board=$(ar71xx_board_name)

case "$board" in
carambola2)
	ucidef_set_system_hostname "Carambola2"
	ucidef_set_wifi_ssid "Carambola2"
	;;

centipede)
	ucidef_set_system_hostname "Centipede"
	ucidef_set_wifi_ssid "Centipede"
	;;

*)
	ucidef_set_system_hostname "OpenWRT"
	;;
esac

uci commit system
uci commit wireless

exit 0
