# target/linux/ipq50xx/base-files/etc/init.d/ath11k_caldata
#!/bin/sh /etc/rc.common

START=10
STOP=90

FIRMWARE="ath11k/IPQ5018/hw1.0/caldata.bin"

ath11kcal_die() {
        echo "ath11cal: " "$*"
        exit 1
}

ath11kcal_extract() {
        local part=$1
        local offset=$2
        local count=$3
        local mtd cal_size

        . /lib/functions.sh
        mtd=$(find_mtd_part $part)
        [ -n "$mtd" ] || \
                ath11kcal_die "no mtd device found for partition $part"

        dd if=$mtd of=/lib/firmware/$FIRMWARE iflag=skip_bytes bs=$count skip=$offset count=1 2>/dev/null || \
                ath11kcal_die "failed to extract calibration data from $mtd"
}

boot() {
        [ -e /lib/firmware/$FIRMWARE ] && exit 0
        ath11kcal_extract "0:ART" 4096 131072
}

