# target/linux/ipq/base-files/etc/uci-defaults/01_leds
#!/bin/sh
#
# Copyright (C) 2011 OpenWrt.org
#

. /lib/functions/uci-defaults.sh
. /lib/ipq806x.sh

board=$(ipq806x_board_name)

case "$board" in
8devices,mango)
	ucidef_set_led_wlan "wlan_2g" "WLAN_2G" "led_2g" "phy1radio"
	ucidef_set_led_wlan "wlan_5g" "WLAN_5G" "led_5g" "phy0radio"
	;;
esac

ucidef_commit_leds

exit 0
